# docker/poller/scripts/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
pidfile=/var/run/supervisord.pid
nodaemon=true
logfile_backups=10  
logfile_maxbytes=50MB  

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:sshd]
command=/usr/sbin/sshd -D -e
priority=10
redirect_stderr=true
stdout_logfile=/var/log/secure
stdout_events_enabled=true
startsecs=5
autorestart=true

[program:centengine]
command=/usr/sbin/centengine /etc/centreon-engine/centengine.cfg
priority=70
process_name=centengine
user=centreon-engine
autorestart=true

[program:snmpd]
command=/usr/sbin/snmpd -f -LS0-6d -Lf /dev/null -p /var/run/snmpd.pid
process_name=snmpd
priority=50
autorestart=true

[program:centreontrapd]
command=/usr/share/centreon/bin/centreontrapd --logfile=/var/log/centreon/centreontrapd.log --severity=error --config=/etc/centreon/conf.pm  --config-extra=/etc/centreon/centreontrapd.pm
priority=70
autorestart=true
user=centreon

[program:snmptrapd]
command=/usr/sbin/snmptrapd -Lsd -f
priority=70
autorestart=true


